# config/mnt/home/Kaba/.config/autostart/welcome.desktop
[Desktop Entry]
Name=Welcome
Exec=/bin/welcome
Type=Application
X-GNOME-Autostart-enabled=true
