# config/inittab
::respawn:/etc/init
::ctrlaltdel:/sbin/reboot
::restart:/sbin/init
